# docker/build/laravel/supervisor.laravel.ini
[unix_http_server]
username = docker
password = docker-password

[supervisord]
logfile = /tmp/supervisord.log
logfile_maxbytes = 50MB
logfile_backups = 10
loglevel = info
pidfile = /tmp/supervisord.pid
nodaemon = true
minfds = 1024
minprocs = 200
user = root

[program:php-fpm]
; command = php-fpm --nodaemonize --fpm-config=/usr/local/etc/php-fpm.d/www.conf
command = php-fpm --nodaemonize
autostart = true
autorestart = true
priority = 5
stdout_logfile = /dev/stdout
stdout_logfile_maxbytes = 0
stderr_logfile = /dev/stderr
stderr_logfile_maxbytes = 0

; https://github.com/laradock/laradock/blob/master/php-worker/supervisord.d/laravel-scheduler.conf.example
[program:laravel-scheduler]
process_name = %(program_name)s_%(process_num)02d
command = /bin/sh -c "while [ true ]; do (php /application/artisan schedule:run --verbose --no-interaction &); sleep 60; done"
autostart = true
autorestart = true
numprocs = 1
stdout_logfile = /dev/stdout
stdout_logfile_maxbytes = 0
stderr_logfile = /dev/stderr
stderr_logfile_maxbytes = 0

; https://github.com/laradock/laradock/blob/master/php-worker/supervisord.d/laravel-worker.conf.example
[program:laravel-worker]
process_name = %(program_name)s_%(process_num)02d
command = php /application/artisan queue:work --sleep=5 --tries=3 --daemon
autostart = true
autorestart = true
numprocs = 1
stdout_logfile = /dev/stdout
stdout_logfile_maxbytes = 0
stderr_logfile = /dev/stderr
stderr_logfile_maxbytes = 0
